# dockerfiles/stravasync.cron
* * * * * root docker exec mycyclediarydotcom_cyclediaryweb_1 python manage.py dbstuff >> /var/log/cron.log 2>&1
0 */1 * * * root docker exec mycyclediarydotcom_cyclediarymongo1_1 mongodump -o /mongobkup >> /var/log/cron.log 2>&1
